# tags/9.0.0.1.29/build/ports/nut/files/upssched-cmd.in
#!/bin/sh
# Copyright (c) 2008-2011 Volker Theile (votdev@gmx.de)
# All rights reserved.

. /etc/rc.subr
. /etc/configxml.subr
. /etc/email.subr

name="upssched-cmd"

load_rc_config "${name}"

# Defaults
nut_upssched_forced_shutdown_timer=${nut_upssched_forced_shutdown_timer:-"0"}

case "${NOTIFYTYPE}" in
	ONLINE)
		_notifymessage="UPS ${UPSNAME} - Running on line power";;
	ONBATT)
		_notifymessage="UPS ${UPSNAME} - Line Power Fail, system running on battery power";;
	LOWBATT)
		_notifymessage="UPS ${UPSNAME} - Low battery power";;
	FSD)
		_notifymessage="UPS ${UPSNAME}: - Forced shutdown in progress";;
	COMMOK)
		_notifymessage="Communications with the UPS ${UPSNAME} are established";;
	COMMBAD)
		_notifymessage="Communications with the UPS ${UPSNAME} are lost";;
	SHUTDOWN)
		_notifymessage="Automatic logout and shutdown proceeding";;
	REPLBATT)
		_notifymessage="UPS ${UPSNAME} - The battery needs to be replaced!";;
	NOCOMM)
		_notifymessage="UPS ${UPSNAME} - The ups is unavailable";;
	NOPARENT)
		_notifymessage="upsmon parent process died - shutdown impossible";;
esac

case "${1}" in
	shutdown-warning)
		_shutdowntimer=`configxml_get "//ups/shutdowntimer"`;
		_message="${_notifymessage}.
Shutdown imminent in ${_shutdowntimer} seconds.";;

	shutdown)
		_message="${_notifymessage}.
Shutdown in progress.";
		shutdown -p now ${_message};;

	resume)
		_message="${_notifymessage}.
Shutdown cancelled.";;

	forced-shutdown)
		_message="${_notifymessage}.
Forced shutdown imminent in ${nut_upssched_forced_shutdown_timer} minutes.";
		shutdown -p +${nut_upssched_forced_shutdown_timer} ${_message};;

	notify)
		_message="${_notifymessage}";;

	*)
		_message="Unknown command: ${1}";
esac

# Write message to syslog
logger -t upssched-cmd "${_message}"

# Send email if feature is enabled
if configxml_isset //ups/email/enable; then
	# Get configured email recipients and subject.
	_recipients=`configxml_get "//ups/email/to" | /usr/bin/tr ";," " "`
	_subject=`configxml_get "//ups/email/subject"`

	# Send email.
	send_email "${_recipients}" "${_subject}" "${_message}"
fi
