# experimental/9.1.0.1/trunk/build/ports/proftpd/files/proftpd.in
#!/bin/sh
#
# Part of NAS4Free (http://www.nas4free.org).
# Copyright (c) 2012-2013 The NAS4Free Project <info@nas4free.org>.
# All rights reserved.
#
# Portions of freenas (http://www.freenas.org).
# Copyright (c) 2005-2011 by Olivier Cochard <olivier@freenas.org>.
# All rights reserved.
#


# PROVIDE: proftpd
# REQUIRE: DAEMON
# BEFORE: LOGIN
# KEYWORD: shutdown
# XQUERY: -i "count(//ftpd/enable) > 0" -o "0" -b
# RCVAR: proftpd

. /etc/rc.subr
. /etc/util.subr
. /etc/configxml.subr

name="proftpd"
rcvar=`set_rcvar`

load_rc_config "$name"

# Custom commands
extra_commands="mkconf"
mkconf_cmd="proftpd_mkconf"
start_precmd="proftpd_mkconf"

# Defaults
proftpd_enable=${proftpd_enable:-"NO"}
proftpd_flags=${proftpd_flags:-""}
proftpd_config=${proftpd_config:-"/var/etc/${name}.conf"}
proftpd_delaytable=${proftpd_delaytable:-"/var/run/${name}/${name}.delay"}
proftpd_displaylogin=${proftpd_displaylogin:-"/var/run/${name}/${name}.motd"}
proftpd_tlsrsacertfile=${proftpd_tlsrsacertfile:-"/var/etc/ssl/private/${name}.cert"}
proftpd_tlsrsakeyfile=${proftpd_tlsrsakeyfile:-"/var/etc/ssl/private/${name}.key"}
proftpd_tcpaccessfilesallow=${proftpd_tcpaccessfilesallow:-"/etc/hosts.allow"}
proftpd_tcpaccessfilesdeny=${proftpd_tcpaccessfilesdeny:-"/etc/hosts.deny"}
proftpd_umaskfile=${proftpd_umaskfile:-"077"}
proftpd_umaskdirectory=${proftpd_umaskdirectory:-"022"}
proftpd_maxinstances=${proftpd_maxinstances:-"none"}
proftpd_defaulttransfermode=${proftpd_defaulttransfermode:-"ascii"}
proftpd_deferwelcome=${proftpd_deferwelcome:-"off"}
proftpd_anonymousrootdir=${proftpd_anonymousrootdir:-"/mnt"}
proftpd_multilinerfc2228=${proftpd_multilinerfc2228:-"off"}
proftpd_timeoutlogin=${proftpd_timeoutlogin:-"300"}
proftpd_timeoutnotransfer=${proftpd_timeoutnotransfer:-"300"}
proftpd_timeoutstalled=${proftpd_timeoutstalled:-"3600"}
proftpd_defaultrootdir=${proftpd_defaultrootdir:-"~"}
proftpd_defaultrootgroupexpr=${proftpd_defaultrootgroupexpr:-"!wheel"}
proftpd_requirevalidshell=${proftpd_requirevalidshell:-"off"}
proftpd_deleteabortedstores=${proftpd_deleteabortedstores:-"off"}
proftpd_tlsverifyclient=${proftpd_tlsverifyclient:-"off"}
proftpd_tlsoptions=${proftpd_tlsoptions:-"NoCertRequest"}
proftpd_tlsprotocol=${proftpd_tlsprotocol:-"SSLv23"}
proftpd_tlsrequired=${proftpd_tlsrequired:-"on"}
proftpd_tcpaccesssysloglevels=${proftpd_tcpaccesssysloglevels:-"info warn"}
proftpd_tcpservicename=${proftpd_tcpservicename:-"ftpd"}
proftpd_timesgmt=${proftpd_timesgmt:-"off"}
proftpd_moddelay_enable=${proftpd_moddelay_enable:-"YES"}
proftpd_modwrap_enable=${proftpd_modwrap_enable:-"YES"}
proftpd_modctrls_controlsengine=${proftpd_modctrls_controlsengine:-"off"}
proftpd_modctrls_controlsacls=${proftpd_modctrls_controlsacls:-"all allow group wheel"}
proftpd_modctrls_controlsauthfreshness=${proftpd_modctrls_controlsauthfreshness:-"10"}
proftpd_modctrls_controlsmaxclients=${proftpd_modctrls_controlsmaxclients:-"5"}
proftpd_modctrls_controlslog=${proftpd_modctrls_controlslog:-"/var/log/${name}/controls.log"}
proftpd_modctrls_controlsinterval=${proftpd_modctrls_controlsinterval:-"10"}
proftpd_modban_bancontrolsacls=${proftpd_modban_bancontrolsacls:-"all allow group wheel"}
proftpd_modban_banlog=${proftpd_modban_banlog:-"/var/log/${name}/ban.log"}
proftpd_modban_banmessage=${proftpd_modban_banmessage:-"Host %a has been banned"}
proftpd_modban_bantable=${proftpd_modban_bantable:-"/var/run/${name}/ban.tab"}
command="/usr/local/sbin/${name}"
command_args="--config ${proftpd_config}"
pidfile="/var/run/${name}.pid"

proftpd_mkconf()
{
	local _welcomemsg _productname _certificate _privatekey _tmpfile

	# Check if required directories exists.
	[ ! -d "/var/run/${name}" ] && mkdir "/var/run/${name}"
	[ ! -d "/var/log/${name}" ] && mkdir "/var/log/${name}"

	# Check if required files exists.
	[ ! -e "${proftpd_delaytable}" ] && cat /dev/null > "${proftpd_delaytable}"
	[ ! -e "${proftpd_tcpaccessfilesallow}" ] && cat /dev/null > "${proftpd_tcpaccessfilesallow}"
	[ ! -e "${proftpd_tcpaccessfilesdeny}" ] && cat /dev/null > "${proftpd_tcpaccessfilesdeny}"
	[ ! -e "/var/log/utx.log" ] && touch /var/log/utx.log && chmod 644 /var/log/utx.log

	# Create welcome message file
	_welcomemsg=`configxml_get "//ftpd/banner"`
	if [ -z "${_welcomemsg}" ]; then
		_productname=`get_product_name`
		echo "Welcome to ${_productname} FTP Server" > ${proftpd_displaylogin}
	else
		echo ${_welcomemsg} > ${proftpd_displaylogin}
	fi

	# Create SSL/TLS certificate
	if configxml_isset //ftpd/tls; then
		_certificate=`configxml_get "//ftpd/certificate"`
		_privatekey=`configxml_get "//ftpd/privatekey"`

		# Create certificate and key files
		if [ -n "${_certificate}" -a -n "${_privatekey}" ]; then
			_tmpfile=/tmp/${name}$$.tmp

			echo "${_certificate}" > ${_tmpfile}
			/usr/bin/uudecode -m -p -r ${_tmpfile} > ${proftpd_tlsrsacertfile}
			/bin/chmod 0600 ${proftpd_tlsrsacertfile}

			echo "${_privatekey}" > ${_tmpfile}
			/usr/bin/uudecode -m -p -r ${_tmpfile} > ${proftpd_tlsrsakeyfile}
			/bin/chmod 0600 ${proftpd_tlsrsakeyfile}

			/bin/rm -f ${_tmpfile}
		fi
	fi

	/usr/local/bin/xml sel -t \
		-v "concat('ServerName \"',//system/hostname,' FTP Server\"')" -n \
		-o "ServerType standalone" -n \
		-o "DefaultServer on" -n \
		-v "concat('DefaultAddress ',//system/hostname)" -n \
		-o "UseIPv6 on" -n \
		-m "//ftpd" \
		-v "concat('Port ',port)" -n \
		-o "User nobody" -n \
		-o "Group nogroup" -n \
		-i "string-length(filemask) > 0 and string-length(directorymask) > 0" -v "concat('Umask ',filemask,' ',directorymask)" -n -b \
		-i "string-length(filemask) > 0 and string-length(directorymask) = 0" -v "concat('Umask ',filemask,' ${proftpd_umaskdirectory}')" -n -b \
		-i "string-length(filemask) = 0 and string-length(directorymask) > 0" -v "concat('Umask ${proftpd_umaskfile} ',directorymask)" -n -b \
		-i "string-length(filemask) = 0 and string-length(directorymask) = 0" -o "Umask ${proftpd_umaskfile} ${proftpd_umaskdirectory}" -n -b \
		-o "SyslogFacility ftp" -n \
		-o "MultilineRFC2228 ${proftpd_multilinerfc2228}" -n \
		-o "DisplayLogin ${proftpd_displaylogin}" -n \
		-o "DeferWelcome ${proftpd_deferwelcome}" -n \
		-i "timeout[. > 0]" -v "concat('TimeoutIdle ',timeout)" -n -b \
		-o "TimeoutLogin ${proftpd_timeoutlogin}" -n \
		-o "TimeoutNoTransfer ${proftpd_timeoutnotransfer}" -n \
		-o "TimeoutStalled ${proftpd_timeoutstalled}" -n \
		-o "MaxInstances ${proftpd_maxinstances}" -n \
		-v "concat('MaxClients ',numberclients)" -n \
		-i "maxconperip[. > 0]" -v "concat('MaxConnectionsPerHost ',maxconperip)" -n -b \
		-i "maxloginattempts[. > 0]" -v "concat('MaxLoginAttempts  ',maxloginattempts )" -n -b \
		-o "DefaultTransferMode ${proftpd_defaulttransfermode}" -n \
		-i "count(fxp) > 0" -o "AllowForeignAddress on" -n -b \
		-i "string-length(pasv_address) > 0" -v "concat('MasqueradeAddress ',pasv_address)" -n -b \
		-i "pasv_min_port[. != '0'] and pasv_max_port[. != '0']" -v "concat('PassivePorts ',pasv_min_port,' ',pasv_max_port)" -n -b \
		-i "count(identlookups) > 0" -o "IdentLookups on" -n -b \
		-i "count(identlookups) = 0" -o "IdentLookups off" -n -b \
		-i "count(usereversedns) > 0" -o "UseReverseDNS on" -n -b \
		-i "count(usereversedns) = 0" -o "UseReverseDNS off" -n -b \
		-i "count(localusersonly) = 0" \
			-n \
			-o "<Anonymous ${proftpd_anonymousrootdir}>" -n \
			-o "  User ftp" -n \
			-o "  Group ftp" -n \
			-o "  UserAlias anonymous ftp" -n \
			-i "string-length(anonymousbandwidth/up) > 0" -v "concat('  TransferRate STOR ',anonymousbandwidth/up)" -n -b \
			-i "string-length(anonymousbandwidth/down) > 0" -v "concat('  TransferRate RETR ',anonymousbandwidth/down)" -n -b \
			-i "count(anonymousonly) > 0" \
				-o "  <Limit LOGIN>" -n \
				-o "    AllowAll" -n \
				-o "  </Limit>" -n \
			-b \
			-o "</Anonymous>" -n \
		-b \
		-i "count(localusersonly) > 0" \
			-n \
			-o "<Limit LOGIN>" -n \
			-o "  AllowGroup ftp" -n \
			-i "count(permitrootlogin) > 0" -o "  AllowGroup wheel" -n -b \
			-o "  DenyAll" -n \
			-o "</Limit>" -n \
		-b \
		-n \
		-o "<Global>" -n \
		-o "  RequireValidShell ${proftpd_requirevalidshell}" -n \
		-i "count(chrooteveryone) > 0" -o "  DefaultRoot ${proftpd_defaultrootdir} ${proftpd_defaultrootgroupexpr}" -n -b \
		-i "count(permitrootlogin) > 0" -o "  RootLogin on" -n -b \
		-o "  AllowOverwrite on" -n \
		-i "count(allowrestart) > 0" \
			-o "  AllowRetrieveRestart on" -n \
			-o "  AllowStoreRestart on" -n \
		-b \
		-o "  DeleteAbortedStores ${proftpd_deleteabortedstores}" -n \
		-i "string-length(userbandwidth/up) > 0" -v "concat('  TransferRate STOR ',userbandwidth/up)" -n -b \
		-i "string-length(userbandwidth/down) > 0" -v "concat('  TransferRate RETR ',userbandwidth/down)" -n -b \
		-o "  TimesGMT ${proftpd_timesgmt}" -n \
		-o "</Global>" -n \
		-i "count(anonymousonly) > 0" \
			-n \
			-o "<Limit LOGIN>" -n \
			-i "count(permitrootlogin) > 0" -o "  AllowGroup wheel" -n -b \
			-o "  DenyAll" -n \
			-o "</Limit>" -n \
		-b \
		-i "count(tls) > 0" \
			-n \
			-o "<IfModule mod_tls.c>" -n \
			-o "  TLSEngine on" -n \
			-o "  TLSProtocol ${proftpd_tlsprotocol}" -n \
			-o "  TLSOptions ${proftpd_tlsoptions}" -n \
			-o "  TLSRSACertificateFile ${proftpd_tlsrsacertfile}" -n \
			-o "  TLSRSACertificateKeyFile ${proftpd_tlsrsakeyfile}" -n \
			-o "  TLSVerifyClient ${proftpd_tlsverifyclient}" -n \
			-i "count(tlsrequired) > 0" \
				-o "  TLSRequired ${proftpd_tlsrequired}" -n \
			-b \
			-o "</IfModule>" -n \
		-b \
		-n \
		-o "<IfModule mod_ban.c>" -n \
		-i "count(mod_ban/enable) = 0" -o "  BanEngine off" -n -b \
		-i "count(mod_ban/enable) > 0" -o "  BanEngine on" -n -b \
		-o "  BanControlsACLs ${proftpd_modban_bancontrolsacls}" -n \
		-o "  BanLog ${proftpd_modban_banlog}" -n \
		-o "  BanMessage ${proftpd_modban_banmessage}" -n \
		-m "mod_ban/rule" \
			-v "concat('  BanOnEvent ',event,' ',occurrence,'/',timeinterval,' ',expire)" -n \
		-b \
		-o "  BanTable ${proftpd_modban_bantable}" -n \
		-o "</IfModule>" -n \
		-n \
		-m "auxparam" \
				-v "." -n \
		-b \
		${configxml_file} | /usr/local/bin/xml unesc > ${proftpd_config}

		# Add mod_delay
		if checkyesno proftpd_moddelay_enable; then
			cat <<EOF >> ${proftpd_config}
<IfModule mod_delay.c>
  DelayEngine on
  DelayTable "${proftpd_delaytable}"
</IfModule>

EOF
		fi

		# Add mod_wrap
		if [ 0 -ne `configxml_get_count "//ftpd/disabletcpwrapper"` ]; then
			proftpd_modwrap_enable="NO"
		fi
		if checkyesno proftpd_modwrap_enable; then
			cat <<EOF >> ${proftpd_config}
<IfModule mod_wrap.c>
  TCPAccessFiles ${proftpd_tcpaccessfilesallow} ${proftpd_tcpaccessfilesdeny}
  TCPAccessSyslogLevels ${proftpd_tcpaccesssysloglevels}
  TCPServiceName ${proftpd_tcpservicename}
</ifModule>

EOF
		fi

		# Add mod_ctrls.
		if [ "${proftpd_modctrls_controlsengine}" = "on" ]; then
			cat <<EOF >> ${proftpd_config}
<IfModule mod_ctrls.c>
  ControlsEngine ${proftpd_modctrls_controlsengine}
  ControlsACLs ${proftpd_modctrls_controlsacls}
  ControlsAuthFreshness
  ControlsMaxClients ${proftpd_modctrls_controlsmaxclients}
  ControlsLog ${proftpd_modctrls_controlslog}
  ControlsInterval ${proftpd_modctrls_controlsinterval}
</IfModule>

EOF
		fi
}

run_rc_command "$1"
