# trunk/build/ports/lcdproc/files/lcdproc.in
#!/bin/sh
# Copyright (C) 2005-2011 FreeNAS Project. All rights reserved.

# PROVIDE: lcdproc
# REQUIRE: LCDd
# XQUERY: -i "count(//nas4free/lcdproc/lcdproc/enable) > 0" -o "0" -b
# RCVAR: lcdproc

. /etc/rc.subr
. /etc/configxml.subr

name="lcdproc"
rcvar=$(set_rcvar)

load_rc_config "${name}"

# Custom commands
start_precmd="lcdproc_mkconf"
start_cmd="lcdproc_start"
stop_cmd="lcdproc_stop"
mkconf_cmd="lcdproc_mkconf"
extra_commands="mkconf"

# Defaults
lcdproc_enable=${lcdproc_enable:-"NO"}
lcdproc_conf=${lcdproc_conf:-"/var/etc/${name}.conf"}

command="/usr/local/bin/${name}"
command_args="-c ${lcdproc_conf}"

lcdproc_mkconf()
{
	# Create config file
	/usr/local/bin/xml sel -t \
		-o '[lcdproc]' -n \
		-o 'Server=localhost' -n \
		-v "concat('Port=', //lcdproc/port)" -n \
		-o 'ReportLevel=2' -n \
		-o 'ReportToSyslog=false' -n \
		-o 'PidFile=/var/run/lcdproc.pid' -n \
		-m "//nas4free/lcdproc/lcdproc/auxparam" \
			-v "." -n \
		-b \
		-n \
		-m "//nas4free/lcdproc/lcdproc/param" \
			-v "." -n \
		-b \
		${configxml_file} | /usr/local/bin/xml unesc > ${lcdproc_conf}
}

lcdproc_start()
{
	local _driver=`configxml_get "//lcdproc/driver"`

	if [ ! "$_rc_prefix" = "force" ]; then
		if [ x"$_driver" = x"curses" -o x"$_driver" = x"text" ]; then
			return 0;
		fi
	fi
	unset start_cmd
	run_rc_command "${_rc_prefix}${rc_arg}" ${rc_extra_args}
}

lcdproc_stop()
{
	local _driver=`configxml_get "//lcdproc/driver"`

	if [ ! "$_rc_prefix" = "force" ]; then
		if [ x"$_driver" = x"curses" -o x"$_driver" = x"text" ]; then
			return 0;
		fi
	fi
	unset stop_cmd
	run_rc_command "${_rc_prefix}${rc_arg}" ${rc_extra_args}
}

run_rc_command "$1"
