# trunk/build/ports/ataidle/files/ataidle.in
#!/bin/sh
#
# Part of NAS4Free (http://www.nas4free.org).
# Copyright (c) 2012-2014 The NAS4Free Project <info@nas4free.org>.
# All rights reserved.
#
# Portions of freenas (http://www.freenas.org).
# Copyright (c) 2005-2011 by Olivier Cochard <olivier@freenas.org>.
# All rights reserved.
#

# PROVIDE: ataidle
# REQUIRE: DAEMON
# BEFORE: LOGIN

. /etc/rc.subr
. /etc/configxml.subr

name="ataidle"

load_rc_config "${name}"

# Defaults
ataidle_delay=${ataidle_delay:-"60"}

echo -n "Initializing Acoustic/Standby/APM on ATA disks: "

# Process all configured IDE disks
/usr/local/bin/xml sel -t -m "//disks/disk[type = 'IDE' or type = 'SCSI']" \
	-v "concat(name,' ',type,' ',devicespecialfile,' ',transfermode,' ',harddiskstandby,' ',apm,' ',acoustic)" \
	-i "position() <= last()" -n -b \
	${configxml_file} | /usr/local/bin/xml unesc | \
	while read _name _type _devicespecialfile _transfermode _harddiskstandby _apm _acoustic; do
		# Skip non ATA disk on CAM
		if [ "SCSI" = "${_type}" -a  "ada" != `/bin/expr -- "${_name}" : '^\([a-zA-Z]*\).*'` ]; then
			continue;
		fi

		# Skip CD/DVD devices
		if [ "acd" = `/bin/expr -- "${_name}" : '^\([a-zA-Z]*\).*'` ]; then
			continue
		fi

		# If UDMA mode is forced, launch atacontrol
		#if [ "auto" != "${_transfermode}" ]; then
		#	/sbin/atacontrol mode ${_name} ${_transfermode} > /dev/null 2>&1
		#fi

		# Set command args
		_cmdargs=""

		if [ ${_acoustic} != 0 ]; then
			_cmdargs="${_cmdargs} -A ${_acoustic}"
		fi

		if [ ${_apm} != 0 ]; then
			_cmdargs="${_cmdargs} -P ${_apm}"
		fi

		if [ ${_harddiskstandby} != 0 ]; then
		        # enter IDLE mode with Standby timer
			_cmdargs="${_cmdargs} -I ${_harddiskstandby}"
			/usr/local/sbin/ataidle ${_cmdargs} ${_devicespecialfile} > /dev/null 2>&1
		else
		    	# enter IDLE mode and disable timer (Always on)
			_cmdargs="${_cmdargs} -I0"
			/usr/local/sbin/ataidle ${_cmdargs} ${_devicespecialfile} > /dev/null 2>&1
		fi

		# Print out some disk information
		echo -n "${_name} "
	done

# Finally issue a line break
echo
