# branches/9.1.0.1/build/ports/inadyn-mt/files/inadyn.in
#!/bin/sh
#
# Part of NAS4Free (http://www.nas4free.org).
# Copyright (c) 2012-2014 The NAS4Free Project <info@nas4free.org>.
# All rights reserved.
# 
# Portions of freenas (http://www.freenas.org).
# Copyright (c) 2005-2011 by Olivier Cochard <olivier@freenas.org>.
# All rights reserved.
#

# PROVIDE: inadyn
# REQUIRE: DAEMON
# KEYWORD: shutdown
# XQUERY: -i "count(//dynamicdns/enable) > 0" -o "0" -b
# RCVAR: inadyn

. /etc/rc.subr
. /etc/configxml.subr

name="inadyn"
rcvar=`set_rcvar`

load_rc_config "${name}"

# Custom commands
start_precmd="inadyn_prestart"

# Defaults
inadyn_enable=${inadyn_enable:-"NO"}
inadyn_config=${inadyn_config:-"/var/etc/${name}.conf"}
command="/usr/local/bin/${name}"
command_args="--input_file ${inadyn_config} > /dev/null"

inadyn_prestart()
{
	local _alias

	/usr/local/bin/xml sel -t -m "//dynamicdns" \
		-i "provider[. = 'dyndns.org']" -o "--dyndns_system dyndns@dyndns.org" -n -b \
		-i "provider[. = 'freedns.afraid.org']" -o "--dyndns_system default@freedns.afraid.org" -n -b \
		-i "provider[. = 'zoneedit.com']" -o "--dyndns_system default@zoneedit.com" -n -b \
		-i "provider[. = 'no-ip.com']" -o "--dyndns_system default@no-ip.com" -n -b \
		-i "provider[. = 'easydns.com']" -o "--dyndns_system default@easydns.com" -n -b \
		-i "provider[. = '3322.org']" -o "--dyndns_system dyndns@3322.org" -n -b \
		-i "provider[. = 'dhis.org']" -o "--dyndns_system default@dhis.org" -n -b \
		-o "--background" -n \
		-o "--syslog" -n \
		-v "concat('--username ',username)" -n \
		-v "concat('--password ',password)" -n \
		-i "string-length(updateperiod) > 0" -v "concat('--update_period_sec ',updateperiod)" -n -b \
		-i "string-length(forcedupdateperiod) > 0" -v "concat('--forced_update_period ',forcedupdateperiod)" -n -b \
		-i "count(wildcard) > 0" -o "--wildcard" -n -b \
		-m "auxparam" \
			-v "." \
			-i "position() <= last()" -n -b \
		-b \
		${configxml_file} | /usr/local/bin/xml unesc > ${inadyn_config}

		for _alias in `configxml_get "//dynamicdns/domainname"`; do
			echo "--alias ${_alias}" >> ${inadyn_config}
		done
}

run_rc_command "$1"
