# trunk/build/ports/msmtp/files/msmtp.in
#!/bin/sh
#
# Part of NAS4Free (http://www.nas4free.org).
# Copyright (c) 2012 The NAS4Free Project <info@nas4free.org>.
# All rights reserved.
# 
# Portions of freenas (http://www.freenas.org).
# Copyright (c) 2005-2011 by Olivier Cochard <olivier@freenas.org>.
# Copyright (c) 2008 Volker Theile (votdev@gmx.de)
# All rights reserved.
#


# PROVIDE: msmtp
# BEFORE: DAEMON

. /etc/rc.subr
. /etc/configxml.subr

name="msmtp"

load_rc_config "${name}"

# Defaults
msmtp_config=${msmtp_config:-"/var/etc/msmtp.conf"}

# Create msmtp.conf configuration
/usr/local/bin/xml sel -t \
	-o "account default" -n \
	-v "concat('host ',//system/email/server)" -n \
	-v "concat('port ',//system/email/port)" -n \
	-o "protocol smtp" -n \
	-v "concat('from ',//system/email/from)" -n \
	-i "count(//system/email/auth) > 0" \
		-v "concat('auth ',//system/email/authmethod)" -n \
		-v "concat('user ',//system/email/username)" -n \
		-v "concat('password ',//system/email/password)" -n \
	-b \
	-i "//system/email/security[. != 'none']" \
		-o "tls on" -n \
		-o "tls_certcheck off" -n \
		-i "//system/email/port[. != '25']" -o "tls_starttls off" -n -b \
		-i "//system/email/port[. = '25']" -o "tls_starttls on" -n -b \
	-b \
	-o "syslog LOG_MAIL" -n \
	${configxml_file} | /usr/local/bin/xml unesc > ${msmtp_config}

# Set user permissions
chmod 0600 ${msmtp_config}
